# scripts/gray-log-monitor.service
[Unit]
Description=Gray Backend Log Monitor
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/gray
EnvironmentFile=-/gray/.env
EnvironmentFile=-/home/ubuntu/gray/.env
ExecStart=/usr/bin/python3 /home/ubuntu/gray/scripts/monitor_logs.py --daemon
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
